# deploy/bita-mock.service
[Unit]
Description=BiTa Mock Backend
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/bita-mock/current
ExecStart=/usr/bin/node /opt/bita-mock/current/server.js
Restart=always
RestartSec=2

# Run as a dedicated non-root user
User=bitasi
Group=bitasi

# Optional overrides (create this file on VPS if needed)
EnvironmentFile=-/etc/bita-mock.env

# Environment
Environment=PORT=8080
Environment=HOST=0.0.0.0
Environment=LOG_REQUESTS=false

# Hardening (reasonable defaults)
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
